# board/raspberrypi5/rootfs_overlay/etc/systemd/system/wpa_supplicant@wlan0.service.d/override.conf
[Unit]
After=systemd-modules-load.service systemd-udevd.service
Wants=systemd-modules-load.service

[Service]
ExecStartPre=/sbin/ip link set %i up
